# crontab
*/1 * * * * python3 /Users/nykon/Development/python/folksonomy/crontab-test.py
